# deploy/fuel-map-bot.service
[Unit]
Description=Fuel Map Telegram Bot
After=network-online.target docker.service
Wants=network-online.target
# Бот обращается к API через host nginx (127.0.0.1:8090)

[Service]
Type=simple
User=root
WorkingDirectory=/opt/fuel-map/telegram-bot
EnvironmentFile=/opt/fuel-map/.env
ExecStart=/usr/bin/node dist/bot.js
Restart=on-failure
RestartSec=10

# Логи: journalctl -u fuel-map-bot -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fuel-map-bot

[Install]
WantedBy=multi-user.target
